# curator4-Ulenk.conf
[program:curator4-Ulenk]
command=python3 /home/curator4/curator4-Ulenk/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-Ulenk/test.err.log
stdout_logfile=/home/curator4/curator4-Ulenk/test.out.log